# service/btsync/run
#!/bin/bash
btsync --config /etc/btsync.conf --nodaemon
