# services/btsync/run
#!/bin/bash

exec /sbin/setuser abc ${APP_PATH}  --config "$CONF_PATH" --nodaemon
