# services/btsync/run
#!/bin/bash

exec /sbin/setuser abc /usr/bin/btsync --config /etc/btsync.conf --nodaemon