# services/btsync/run
#!/bin/bash

exec /sbin/setuser abc /app/btsync --config /config/btsync.conf --nodaemon